# packaging/systemd/latch-cron-hourly.timer
[Unit]
Description=Run Latch hourly maintenance

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target